# custom_components/skelly_ultra/skelly_ultra_srv/skelly-ultra-server.service
[Unit]
Description=Skelly Ultra REST Server for Bluetooth Classic Audio
Documentation=https://github.com/martinecker/home-assistant-core
After=network.target bluetooth.target pipewire.service
Wants=bluetooth.target

[Service]
Type=simple
# Run as root to enable automatic Bluetooth pairing via D-Bus
# If you don't need automatic pairing, you can set User=YOUR_USERNAME_HERE instead
User=root
Group=root
WorkingDirectory=/path/to/custom_components/skelly_ultra
ExecStart=/usr/bin/python3 -m skelly_ultra_srv.server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
# Note: NoNewPrivileges and ProtectSystem are relaxed to allow D-Bus system bus access
PrivateTmp=true
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
